# trash-empty.service
[Unit]
Description=daily trash-empty service

[Service]
ExecStart=/usr/bin/trash-empty
ExecStopPost=/bin/echo 'Finished trash empty'

